# dist/etc/init.d/krglast
#! /bin/sh
### BEGIN INIT INFO
# Provides:          krglast
# Required-Start:    $all
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Description:       Disable rcS in Kerrighed container
### END INIT INFO

# Author: Jean Parpaillon <jean.parpaillon@kerlabs.com>
#

# PATH should only include /usr/* if it runs after the mountnfs.sh script
PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="Enabling specific rcS for Kerrighed container"
NAME=krglast

. /lib/init/kerrighed-functions

is_kerrighed_kernel || exit 0
is_kerrighed_boot || exit 0

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.2-14) to ensure that this file is present
# and status_of_proc is working.
. /lib/lsb/init-functions

#
# Function that starts the daemon/service
#
do_start()
{
    [ -d /etc/rcS.d.krg ] && rm -r /etc/rcS.d.krg
    mkdir -p /etc/rcS.d.krg

    # Replace original rcS with Kerrighed container's one
    mount -o bind /etc/rcS.d.krg /etc/rcS.d
    RET=$?

    return $RET
}

#
# Function that stops the daemon/service
#
do_stop()
{
    return 0
}

#
# Function that sends a SIGHUP to the daemon/service
#
do_reload() {
	return 0
}

case "$1" in
  start)
	log_daemon_msg "$DESC"
	do_start
	case "$?" in
		0|1) log_end_msg 0 ;;
		2) log_end_msg 1 ;;
	esac
	;;
  stop)
	log_daemon_msg "Stopping $NAME"
	do_stop
	case "$?" in
		0|1) log_end_msg 0 ;;
		2) log_end_msg 1 ;;
	esac
	;;
  *)
	echo "Usage: $SCRIPTNAME {start|stop}" >&2
	exit 3
	;;
esac

:
