# dist/etc/init.d/krgreboot
#! /bin/sh
### BEGIN INIT INFO
# Provides:          krgreboot
# Required-Start:
# Required-Stop:     reboot
# Should-Stop:
# X-Stop-After:      umountroot
# Default-Start:
# Default-Stop:      6
# Short-Description: Execute the reboot command.
# Description:
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin

. /lib/lsb/init-functions
. /lib/init/kerrighed-functions

is_kerrighed_kernel || exit 0

do_stop () {
    if ! is_kerrighed_up; then
	log_warning_msg "Kerrighed not running: doing normal reboot"
	return 0
    fi	

	log_action_msg "Will now restart the cluster"
	krgadm cluster reboot

	# Do not return to not run the regular halt
	while true; do
	    sleep 1
	done
}

case "$1" in
  start)
	# No-op
	;;
  restart|reload|force-reload)
	echo "Error: argument '$1' not supported" >&2
	exit 3
	;;
  stop)
	do_stop
	;;
  *)
	echo "Usage: $0 start|stop" >&2
	exit 3
	;;
esac
